# files/onedrive/media-onedrive.automount
Description=Automount Onedrive

[Automount]
Where=/media/onedrive

[Install]
WantedBy=multi-user.target
